# roles/common-wall/files/systemd/kissvpn-client@.service
[Unit]
Description=Kissvpn client
After=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/kissvpn/$i.conf
ExecStart=/usr/local/bin/kissvpn --key @/etc/kissvpn/$i.key --up_script /etc/kissvpn/$i.up.sh connect ${KISSVPN_REMOTE}
Restart=always
RestartSec=30
StartLimitBurst=10

[Install]
WantedBy=multi-user.target
